# motd-scripts/update-motd.d/00-header
#!/bin/sh
#
# 00-header
# OS-Infos laden
[ -r /etc/os-release ] && . /etc/os-release
OS=$PRETTY_NAME

# Fallback, falls PRETTY_NAME leer ist
if [ -z "$OS" ] && [ -x /usr/bin/lsb_release ]; then
    OS="$(lsb_release -s -d)"
fi

# Pfad zu lolcat (mit --force, da pam_motd kein TTY ist)
LOLCAT="/usr/local/bin/lolcat -f"

# Hostname im FIGlet-Slant-Font
if [ -x /usr/local/bin/lolcat ]; then
    figlet -f slant "$(hostname)" | /usr/local/bin/lolcat -f
else
    figlet -f slant "$(hostname)"
fi

printf "\n"

# Info-Block bauen
INFO=$(printf "\t- %s\n\t- OS version %s\n\t- Kernel %s\n" \
  "$OS" "$(cat /etc/debian_version 2>/dev/null || echo "n/a")" "$(uname -r)")

# Info-Block farbig (oder plain als Fallback)
if [ -x /usr/local/bin/lolcat ]; then
    printf "%s\n" "$INFO" | /usr/local/bin/lolcat -f
else
    printf "%s\n" "$INFO"
fi

printf "\n"
